# packages/helpers/woc-api/woc-api.service
# Place this file in /etc/systemd/system
# systemctl start woc-api
# systemctl stop woc-ap
# journalctl -xe
# systemctl enable woc-api

[Unit]
Description=WoC API Server
After=network.target

[Service]
LimitNOFILE=infinity
LimitNPROC=infinity
LimitCORE=infinity
ExecStart=/home/ubuntu/woc-api/woc-api
WorkingDirectory=/home/ubuntu/woc-api
User=ubuntu
Environment="SETTINGS_CONTEXT=live.<woc|satoshi>.<name>.<network>"
Type=simple
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
